# supervisord.conf
[program:celeryd]
command={{ PYTHON }} {{ PROJECT_DIR }}/manage.py celeryd

[program:celerybeat]
command={{ PYTHON }} {{ PROJECT_DIR }}/manage.py celerybeat
